# mcp251xfd/99-devcoredump.rules
ACTION=="add", SUBSYSTEM=="devcoredump", RUN+="/usr/sbin/devcoredump"
